# Request Mirror/files/request-mirror.service
[Unit]
Description=Request Mirror

[Service]
User=requestmirror
ExecStart=/usr/local/bin/request-mirror

[Install]
WantedBy=multi-user.target